# workspace/self_perception_simple_clean/scripts/launchd/com.self-perception.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.self-perception.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/Administrator/.openclaw/workspace/self_perception_simple_clean/watchdog.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/Administrator/.openclaw/workspace/self_perception_simple_clean</string>
    <key>StandardOutPath</key>
    <string>/Users/Administrator/.openclaw/workspace/self_perception_simple_clean/data/logs/watchdog_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Administrator/.openclaw/workspace/self_perception_simple_clean/data/logs/watchdog_stderr.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <!-- watchdog 依赖 supervisor，先启动 supervisor -->
    <key>KeepAlive</key>
    <dict>
        <key>OtherJobEnabled</key>
        <string>com.self-perception.supervisor</string>
        <true/>
    </dict>
</dict>
</plist>
